# root/etc/services.d/sinusbot/run
#!/usr/bin/with-contenv bash

touch /config/sinusbot.log

s6-setuidgid abc /config/sinusbot/sinusbot >> /config/sinusbot.log
